# uuid.service
[Unit]
Description=Golang Webserver that generates UUID on :/host_id

[Service]
ExecStart=/home/ec2-user/main

[Install]
WantedBy=multi-user.target
